# packages/base/src/files/etc/cron.d/audit.cron
0 0 * * * . ${{POLICY_HOME}}/etc/profile.d/env.sh; elk.sh --audit < /dev/null > /dev/null 2>&1
