# packages/base/src/files/etc/cron.d/monitor.cron
* * * * * . ${{POLICY_HOME}}/etc/profile.d/env.sh; monitor.sh </dev/null > /dev/null 2>&1
